# core/docker/defaults/supervisor/config/conf.d/minemeld-traced.conf
[program:minemeld-traced]
environment=PATH="/opt/minemeld/engine/current/bin",PYTHONPATH="/opt/minemeld/engine",HOME="/opt/minemeld/local/trace"
command=/opt/minemeld/engine/current/bin/mm-traced /opt/minemeld/local/config/traced/traced.yml
directory=/opt/minemeld/local/trace
umask=027
stopwaitsecs=200
autostart=true
autorestart=true
redirect_stderr=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
